# stage-bugg-os/03-automount-sd-service/mnt-sd.automount
[Unit]
Description=Automount SD Card

[Automount]
Where=/mnt/sd
TimeoutIdleSec=0

[Install]
WantedBy=multi-user.target
